# roles/kafka/templates/kafka.service.j2
[Unit]
Description=Apache Kafka
Documentation=https://kafka.apache.org/documentation/
After=network-online.target zookeeper.service
Wants=network-online.target
Requires=network-online.target

[Service]
Type=simple
User={{ kafka_user }}
Group={{ kafka_group }}
Environment="JAVA_HOME={{ java_home | default('/usr/lib/jvm/jre') }}"
Environment="LOG_DIR={{ kafka_log_dir }}"
EnvironmentFile=-{{ kafka_conf_dir }}/kafka-env.sh
ExecStart={{ kafka_install_dir }}/bin/kafka-server-start.sh {{ kafka_config_file }}
ExecStop={{ kafka_install_dir }}/bin/kafka-server-stop.sh
Restart=on-failure
RestartSec=10
TimeoutStartSec=300
TimeoutStopSec=120

# Resource limits
LimitNOFILE=65536
LimitNPROC=32768

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths={{ kafka_data_dir }}
ReadWritePaths={{ kafka_log_dir }}
ReadWritePaths={{ kafka_log_dirs }}

[Install]
WantedBy=multi-user.target
